# root/etc/uci-defaults/luci-app-wolh
#!/bin/sh

[ ! -f /etc/config/wolh ] && touch /etc/config/wolh